# systemd/composed_camera_server.service
# systemd service for the SONIC composed camera server.
#
# Publishes JPEG-encoded camera frames over ZMQ so the data exporter
# (on the workstation) can subscribe and record.
#
# ── Setup ────────────────────────────────────────────────────────
#  The install script can generate and install this automatically:
#       bash install_scripts/install_camera_server.sh
#
#  Or manually:
#  1. Edit USER, HOME, REPO_DIR, and the ExecStart command below.
#  2. Copy to systemd:
#       sudo cp systemd/composed_camera_server.service /etc/systemd/system/
#  3. Enable and start:
#       sudo systemctl daemon-reload
#       sudo systemctl enable composed_camera_server.service
#       sudo systemctl start composed_camera_server.service
#  4. Check status:
#       sudo systemctl status composed_camera_server.service
#       journalctl -u composed_camera_server.service -f
#
# ── Finding your device ID ───────────────────────────────────────
#  OAK cameras:
#       python -c "import depthai as dai; print(dai.Device.getAllAvailableDevices())"
#  RealSense:
#       rs-enumerate-devices --short
#  USB:
#       ls /dev/video*    # use the index number, e.g. "0" for /dev/video0
# ─────────────────────────────────────────────────────────────────

[Unit]
Description=SONIC Composed Camera Server (ZMQ)
After=network.target

[Service]
Type=simple

# ── EDIT THESE for your robot ────────────────────────────────────
User=unitree
Environment="HOME=/home/unitree"
Environment="REPO_DIR=/home/unitree/Humanoid/GR00T-WholeBodyControl"
# ─────────────────────────────────────────────────────────────────

WorkingDirectory=/home/unitree/Humanoid/GR00T-WholeBodyControl

# Unitree G1 wrist cameras (Intel RealSense D405)
# SN=260522278617 -> left wrist (index 0), SN=260622271155 -> right wrist (index 1)
ExecStart=/home/unitree/Humanoid/GR00T-WholeBodyControl/.venv_camera/bin/python -m gear_sonic.camera.composed_camera \
    --left-wrist-camera realsense --left-wrist-device-id 0 \
    --right-wrist-camera realsense --right-wrist-device-id 1 \
    --port 5555

# Full 3-camera setup (add when head camera is resolved):
# ExecStart=/home/unitree/Humanoid/GR00T-WholeBodyControl/.venv_camera/bin/python -m gear_sonic.camera.composed_camera \
#     --ego-view-camera unitree_g1 --ego-view-device-id enP8p1s0 \
#     --left-wrist-camera realsense --left-wrist-device-id 0 \
#     --right-wrist-camera realsense --right-wrist-device-id 1 \
#     --port 5555

Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
